# Raspberry Pi/Raspi-Setup-onboard/autostart/keyboard.desktop
[Desktop Entry]
Name=Screen Keyboard Autostart
Exec=/home/pi/Python/keyboard.sh
Type=Application
Terminal=false

