# Raspberry Pi/Raspi-Setup/autostart/auto-copy-files-from-usb.desktop
[Desktop Entry]
Name=rob file updating
Exec=/home/pi/Python/auto-copy-files-from-usb.sh
Type=Application
Terminal=false