# Raspberry Pi/Raspi-Setup/autostart/keyboard-resize.desktop
[Desktop Entry]
Name=Move Screen Keyboard
Exec=/home/pi/Python/move-keyboard.sh
Type=Application
Terminal=false
