# Raspberry Pi/Raspi-Setup/autostart/Multipack_Parser.desktop
[Desktop Entry]
Name=Multipack Parser Autostart
Exec=/home/pi/Python/multipack-parser.sh
Type=Application
Terminal=false
